# deploy/supervisor/personal_celestialmc_pw-schedule.conf
[program:personal_celestialmc_pw-schedule]
process_name=%(program_name)s
command=/www/server/php/83/bin/php artisan schedule:work
directory=/www/wwwroot/personal.celestialmc.pw
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/www/wwwroot/personal.celestialmc.pw/storage/logs/schedule-worker.log